# systemd/mqtt_logger.service
[Unit]
Description=MQTT Logger (persists readings to SQLite + publishes alerts)
Wants=network-online.target
After=network-online.target

[Service]
User=rs
Group=rs
WorkingDirectory=/home/rs/weather
ExecStart=/home/rs/weather/venv/bin/python3 /home/rs/weather/mqtt_logger.py
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
